# invoiceninja/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
exec php-fpm7 --nodaemonize