# invoiceninja/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/bin/bash
# Start php-fpm in the foreground
exec php-fpm -F