# invoiceninja/rootfs/etc/s6-overlay/s6-rc.d/ninja/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InvoiceNinja
# # Start the original docker-entrypoint and supervisord under s6
# ==============================================================================

# Log the values for debugging
bashio::log.info "Starting Invoice Ninja with APP_URL=${APP_URL}"
bashio::log.info "Database host: ${DB_HOST}"

exec /usr/local/bin/docker-entrypoint -n
